# linux/.config/polybar/scripts/system-usb-udev/95-usb.rules
KERNEL=="sd*", ACTION=="add", ATTR{removable}=="1", \
    RUN+="/home/ishaan/.config/polybar/scripts/system-usb-udev/system-usb-udev.sh --update"
KERNEL=="sd*", ACTION=="remove", \
    RUN+="/home/ishaan/.config/polybar/scripts/system-usb-udev/system-usb-udev.sh --update"
